# daemon/geiger.service
[Unit]
Description=Geiger daemon

[Service]
ExecStart=/usr/local/bin/geiger

[Install]
WantedBy=default.target
